# dist/systemd/system/greenboot-rollback.service

[Unit]
Description=Greenboot rollback
DefaultDependencies=no
Before=multi-user.target
After=greenboot.service ostree-finalize-staged.service
Requires=ostree-finalize-staged.service
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no

[Install]
WantedBy=multi-user.target