# dist/systemd/system/greenboot-rollback.service
[Unit]
Description=Greenboot rollback
DefaultDependencies=no
Before=systemd-update-done.service
Wants=local-fs.target 
After=local-fs.target
ConditionNeedsUpdate=|/var
ConditionNeedsUpdate=|/etc

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no

[Install]
WantedBy=multi-user.target