# dist/systemd/system/greenboot-rollback.service
[Unit]
Description=Greenboot rollback
Requires=dbus.service
RequiresMountsFor=/boot
ConditionNeedsUpdate=|/var
ConditionNeedsUpdate=|/etc

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no