# dist/systemd/system/greenboot.service
[Unit]
Description=greenboot Health Checks Runner
Before=boot-complete.target
OnFailureJobMode=fail
RequiresMountsFor=/boot
RequiresMountsFor=/etc
OnFailure=greenboot-rollback.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/greenboot/greenboot health-check
Restart=no

[Install]
RequiredBy=boot-complete.target
WantedBy=multi-user.target