# dist/systemd/system/greenboot.service
[Unit]
Description=greenboot Health Checks Runner
DefaultDependencies=no
Before=boot-complete.target greenboot-rollback.service
OnFailureJobMode=fail

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/greenboot/greenboot health-check
Restart=no

[Install]
RequiredBy=boot-complete.target
WantedBy=multi-user.target