# dist/systemd/system/greenboot-rollback.service
[Unit]
Description=Greenboot rollback
Requires=systemd-update-done.service
After=local-fs.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no