# dist/systemd/system/greenboot-rollback.service
[Unit]
Description=Greenboot rollback
Wants=systemd-update-done.service
After=local-fs.target systemd-update-done.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no