# dist/systemd/system/greenboot.service
[Unit]
Description=greenboot Health Checks Runner
Before=boot-complete.target greenboot-rollback.service systemd-update-done.service
OnFailureJobMode=fail

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/greenboot/greenboot health-check
Restart=no

[Install]
RequiredBy=boot-complete.target
WantedBy=multi-user.target