# dist/systemd/system/greenboot-rollback.service
[Unit]
Description=Greenboot rollback
Requires=dbus.service
RequiresMountsFor=/boot
Wants=systemd-update-done.service
After=systemd-update-done.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot rollback
Restart=no